# supervisor/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9001

[supervisord]
logfile=/var/log/supervisord.log
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
serverurl=http://127.0.0.1:9001

[program:mongo]
command=mongod --config /usr/local/etc/mongod.conf --rest
stdout_logfile=/var/log/supervisor-mongod.log
stderr_logfile=/var/log/supervisor-mongod.err.log
autostart=false

[program:socket]
command=nohup node socket/index.js &
environment=PORT='1347'
stdout_logfile=/var/log/supervisor-socket.log
stderr_logfile=/var/log/supervisor-socket.err.log
autostart=false

[program:api]
command=nohup api &
environment=PORT='3000',MONGOLAB_URI='mongodb://localhost/heroku_app37083199'
stdout_logfile=/var/log/supervisor-api.log
stderr_logfile=/var/log/supervisor-api.err.log
autostart=false

[program:web]
command=python -m SimpleHTTPServer 8000
directory=web/public
stdout_logfile=/var/log/supervisor-web.log
stderr_logfile=/var/log/supervisor-web.err.log
autostart=false
