# supervisor/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9001

[supervisord]
logfile=/var/log/ftdd.log
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
serverurl=http://127.0.0.1:9001

[program:mongo]
command=mongod --config /usr/local/etc/mongod.conf --rest
stdout_logfile=/var/log/ftd-mongod.log
stderr_logfile=/var/log/ftd-mongod.err.log
priority=1

[program:redis]
command=/usr/local/bin/redis-server
stdout_logfile=/var/log/ftd-%(program_name)s.log
stderr_logfile=/var/log/ftd--%(program_name)s.err.log
priority=1

[program:socket]
command=nohup node socket/index.js &
environment=PORT='1347'
stdout_logfile=/var/log/ftd-%(program_name)s.log
stderr_logfile=/var/log/ftd-%(program_name)s.err.log
priority=2
autostart=false

[program:drakov]
command=nohup node api/index.js &
stdout_logfile=/var/log/ftd-%(program_name)s.log
stderr_logfile=/var/log/ftd-%(program_name)s.err.log
priority=2
autostart=false

[program:api]
command=nohup api &
environment=PORT='3000',MONGOLAB_URI='mongodb://localhost/heroku_app37083199',DB='heroku_app37083199'
stdout_logfile=/var/log/ftd-%(program_name)s.log
stderr_logfile=/var/log/ftd-%(program_name)s.err.log
priority=2
autostart=false

[program:web]
command=python -m SimpleHTTPServer 8000
directory=web/public
stdout_logfile=/var/log/ftd-%(program_name)s.log
stderr_logfile=/var/log/ftd-%(program_name)s.err.log
priority=2
autostart=false
