# supervisor/supervisord.conf
[unix_http_server]
file=tmp/supervisor.sock   ; (the path to the socket file)

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9001

[supervisord]
logfile=logs/supervisord.log
pidfile=tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
serverurl=http://127.0.0.1:9001

[program:socket]
command=nohup node socket/index.js &
stdout_logfile=logs/supervisor-socket.log
stderr_logfile=logs/supervisor-socket.err.log
autostart=false

[program:api]
command=nohup node api/index.js &
stdout_logfile=logs/supervisor-api.log
stderr_logfile=logs/supervisor-api.err.log
autostart=false

[program:web]
command=python -m SimpleHTTPServer 8000
directory=web/public
stdout_logfile=logs/supervisor-web.log
stderr_logfile=logs/supervisor-web.err.log
autostart=false
